# build/mysqlgd.service
[Unit]
Description=MySQL Guardian Monitoring Service
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/mysql-guardian --demonize
Restart=always
WorkingDirectory=/var/log
User=root
Group=root

[Install]
WantedBy=multi-user.target
